# scripts/crontab
MAILTO=helm@mowgli.cs.unibo.it
HOME=/home/tassi/
#SVNOPTIONS='-r {2006-01-09}'
00 02 * * * sh /home/tassi/helm/software/matita/scripts/crontab.sh
